# .config/systemd/user/hawki-openai-proxy.service
[Unit]
Description=Local OpenAI-compatible proxy for HAWKI
Documentation=file:%h/.local/share/hawki-openai-proxy/README.md
After=network-online.target
Wants=network-online.target
# HAWKI is only reachable through the SSH SOCKS tunnel. Wants= (not Requires=)
# on purpose: if the tunnel is down, the proxy must stay up and answer 504 with
# a clear message rather than disappear and take /healthz with it.
Wants=hawki-socks-tunnel.service
After=hawki-socks-tunnel.service

[Service]
Type=simple
WorkingDirectory=%h/.local/share/hawki-openai-proxy
# --host is hardcoded to 127.0.0.1: the proxy is unauthenticated and must never
# be reachable from the network. --no-access-log keeps request URLs out of the
# journal; the app logs status/model/duration/category itself.
ExecStart=%h/.local/share/hawki-openai-proxy/.venv/bin/uvicorn hawki_proxy.app:app --host 127.0.0.1 --port 11435 --no-access-log --log-level warning
Restart=on-failure
RestartSec=5

# The token is deliberately absent here. It is read from KeePassXC per request.
# No Environment=, no EnvironmentFile=, no LoadCredential=.
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
CapabilityBoundingSet=

[Install]
WantedBy=default.target
