# .config/systemd/user/waybar.service.d/log-filter.conf
[Service]
# Suppress known-harmless, high-frequency waybar tray noise:
#   - "... is already registered"  (Dropbox/nm-applet SNI re-registration)
#   - "GTK_IS_ACCEL_GROUP ... failed" (GTK3/libdbusmenu accel assertion, upstream bug)
#
# LogFilterPatterns cannot catch the GTK lines: GLib sends them via the NATIVE journal
# protocol (_TRANSPORT=journal), which bypasses stream-based LogFilterPatterns. Routing
# waybar's output through a grep pipe makes stderr a pipe (not the journal socket), so
# GLib falls back to plain-text logging and BOTH message classes become filterable in one
# place. Restart=always in the base unit still restarts waybar if it exits.
ExecStart=
ExecStart=/bin/bash -c 'exec /usr/bin/waybar 2>&1 | grep --line-buffered -vE "GTK_IS_ACCEL_GROUP|is already registered"'
